# docker/supervisor/conf.d/messenger-ai.conf
[program:messenger-ai]
command=php /var/www/html/app/bin/console messenger:consume ai_deepseek_jobs ai_reply_send --memory-limit=256M
process_name=%(program_name)s_%(process_num)02d
numprocs=3
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
